# scripts/config/deploy-dev-supervisor-api.conf
[program:ichrisbirch-api]
directory=/usr/local/var/www/ichrisbirch/ichrisbirch
command=/usr/local/var/www/ichrisbirch/.venv/bin/gunicorn wsgi:api --bind 0.0.0.0:4200 --workers 4 --worker-class uvicorn.workers.UvicornWorker
environment=ENVIRONMENT="development"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/usr/local/var/log/supervisor/ichrisbirch-api.log
stout_logfile=/usr/local/var/log/supervisor/ichrisbirch-api.log
