# deploy/prod/supervisor-app.conf
[program:euphoria-app]
directory=/var/www/euphoria/euphoria
command=/var/www/euphoria/.venv/bin/gunicorn wsgi:app --bind 0.0.0.0:8000 --workers 4
environment=ENVIRONMENT="production"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/euphoria-app-error.log
stout_logfile=/var/log/supervisor/euphoria-app-out.log